# templates/auto-ntpdate.timer.j2
[Unit]
Description=ANXS auto NTP date sync timer

[Timer]
OnCalendar=*-*-* 00/{{ ntp_auto_update_hour }}:00:00
Persistent=true

[Install]
WantedBy=timers.target
